# scripts/systemd/first-boot.service
[Unit]
Description=Sugarkube first boot verifier
After=network-online.target cloud-final.service
Wants=network-online.target cloud-final.service
ConditionPathExists=!/var/log/sugarkube/first-boot.ok

[Service]
Type=oneshot
ExecStart=/opt/sugarkube/first_boot_service.py
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
