# scripts/systemd/ssd-clone.service
[Unit]
Description=Sugarkube SSD clone automation
After=first-boot.service systemd-udevd.service
Wants=first-boot.service
ConditionPathExists=/opt/sugarkube/ssd_clone.py
ConditionPathExists=!/var/log/sugarkube/ssd-clone.done

[Service]
Type=oneshot
EnvironmentFile=-/etc/sugarkube/notifications.env
ExecStartPre=/bin/udevadm settle --timeout=30
ExecStart=/opt/sugarkube/ssd_clone_service.py
StandardOutput=journal
StandardError=journal
