# scripts/systemd/sugarkube-export-kubeconfig.path
[Unit]
Description=Trigger kubeconfig export when k3s kubeconfig changes

[Path]
PathExists=/etc/rancher/k3s/k3s.yaml
PathModified=/etc/rancher/k3s/k3s.yaml
Unit=sugarkube-export-kubeconfig.service

[Install]
WantedBy=multi-user.target
