# scripts/systemd/sugarkube-export-kubeconfig.service
[Unit]
Description=Export sanitized kubeconfig to /boot for offline retrieval
After=k3s.service
Wants=k3s.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/sugarkube-export-kubeconfig.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
