# scripts/systemd/sugarkube-export-kubeconfig.service
[Unit]
Description=Export sanitized kubeconfig to /boot
After=k3s.service
Requires=k3s.service
ConditionPathIsDirectory=/boot

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/sugarkube-export-kubeconfig.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
