# scripts/systemd/sugarkube-export-kubeconfig.service
[Unit]
Description=Export the Sugarkube kubeconfig to the boot partition
After=network-online.target
Wants=network-online.target
ConditionPathExists=/opt/sugarkube/export-kubeconfig.sh

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/opt/sugarkube/export-kubeconfig.sh

[Install]
WantedBy=multi-user.target
